# app/conf.d/unity.conf
[program:unity]
command=startkde
autorestart=true
