# app/conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :0 -screen 0 1280x800x24
autorestart=true
